# roles/nginx/files/certbot.service
[Unit]
Description=Renew nginx certificate

[Service]
Type=oneshot
ExecStart=/opt/certbot/bin/certbot renew -q
