# roles/spy/files/spy.service
[Unit]
Description=Spy

[Service]
Type=exec
ExecStart=/usr/bin/spy

[Install]
WantedBy=multi-user.target
